# ubiquiti/conf/etc/init/unifi-exporter.conf
start on started networking and runlevel [2345]
stop on runlevel [016]
respawn
setuid unifi
setgid unifi
exec {{ server_bin }} --config.file="{{ config_file }}" 
